# ansible/roles/postgresql/templates/pg-backup.timer.j2
# /etc/systemd/system/pg-backup.timer

[Unit]
Description=PostgreSQL Backup
Requires=pg-backup.service

[Timer]
#Unit==pg-backup.service
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
